# deploy/deployment-0/com.runforge.daemon0.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.runforge.daemon0</string>
  <!-- Deployment #0: governed self-hosted runforge daemon.
       Runs launch-deployment0.sh (runtime clone ~/code/runforge-runtime,
       DB runforge_prod0 on the compose postgres :45432, decision index enabled).
       Distinct label from the legacy com.runforge.daemon so scripts/install-daemon.sh
       never clobbers this. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>~/code/runforge-ops/launch-deployment0.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:~/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>~</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>~/code/runforge-ops/logs/daemon0.log</string>
  <key>StandardErrorPath</key>
  <string>~/code/runforge-ops/logs/daemon0.log</string>
</dict>
</plist>
