# scripts/com.autoclaude.reviewer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.autoclaude.reviewer</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>~/code/auto-claude/scripts/reviewer.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>~/code/auto-claude</string>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>~/logs/claude-reviewer.log</string>
  <key>StandardErrorPath</key>
  <string>~/logs/claude-reviewer.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>~</string>
  </dict>
</dict>
</plist>
