# scripts/com.runforge.docker-autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.runforge.docker-autostart</string>
  <!-- Launch Docker Desktop (in the background) at login so the pilot daemon
       container (restart:unless-stopped) resurrects after a reboot. One-shot:
       `open` returns immediately and Docker starts its VM/engine asynchronously;
       the container then auto-restarts once the engine is up. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/open</string>
    <string>-g</string>
    <string>-a</string>
    <string>Docker</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/logs/docker-autostart.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/logs/docker-autostart.log</string>
</dict>
</plist>
