# scripts/com.runforge.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.runforge.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NPX_PATH__</string>
    <string>tsx</string>
    <string>packages/daemon/src/main.ts</string>
    <string>start</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>__HOME__/logs/claude-daemon.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/logs/claude-daemon.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>GITHUB_TOKEN</key>
    <string>__GITHUB_TOKEN__</string>
    <key>RUNFORGE_DATABASE_URL</key>
    <string>__RUNFORGE_DATABASE_URL__</string>
    <key>DAEMON_DATA_BACKEND</key>
    <string>__DAEMON_DATA_BACKEND__</string>
    <key>ENCRYPTION_KEY</key>
    <string>__ENCRYPTION_KEY__</string>
    <key>RUNFORGE_CONTROL_TOKEN</key>
    <string>__RUNFORGE_CONTROL_TOKEN__</string>
  </dict>
</dict>
</plist>
