# scripts/com.autoclaude.creds-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.autoclaude.creds-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/scripts/sync-claude-creds.sh</string>
    <string>__CREDS_DIR__</string>
  </array>
  <!-- Refresh the container's creds well inside the access-token validity so
       the container never needs to run its own (refresh-token-rotating) refresh. -->
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>__HOME__/logs/claude-creds-sync.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/logs/claude-creds-sync.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>
</dict>
</plist>
